# infra.xml
<?xml version='1.0'?>
<!DOCTYPE service_bundle SYSTEM '/usr/share/lib/xml/dtd/service_bundle.dtd.1'>

<!--   
         Copyright 2006-2010 Yuri Voinov.   All rights reserved.
         Service manifest for Oracle iAS Infrastructure services.
-->

<service_bundle type='manifest' name='oracle-services'>

<service
        name='application/ias-infra'
        type='service'
        version='1'>

        <create_default_instance enabled='false' />

        <single_instance />

        <dependency name="multi-user" 
                grouping="require_all" 
              restart_on="error" 
                    type="service">
           <service_fmri value="svc:/milestone/multi-user:default"/>
        </dependency>

        <exec_method
                type='method'
                name='start'
                exec='/lib/svc/method/init.ias_infra %m'
                timeout_seconds='0' />

        <exec_method
                type='method'
                name='stop'
                exec='/lib/svc/method/init.ias_infra %m'
                timeout_seconds='0' />

        <exec_method
                type='method'
                name='restart'
                exec='/lib/svc/method/init.ias_infra %m'
                timeout_seconds='0' />

        <property_group name='options' type='application'>
                <propval name='tcp_listen' type='boolean'
                        value='true' />
        </property_group>

        <stability value='External' />

        <template>
                <common_name>
                        <loctext xml:lang='C'> 
                         Oracle iAS Infrastructure Services
                        </loctext>
                </common_name>
        </template>
</service>

</service_bundle>